# backend/deploy/aifficient-celery-beat.service
[Unit]
Description=Aifficient Celery Beat (campaign scheduler tick)
After=network.target redis.service aifficient-celery-worker.service
Wants=redis.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Afficient-ai/backend
EnvironmentFile=/home/ubuntu/Afficient-ai/backend/.env
ExecStart=/home/ubuntu/Afficient-ai/backend/venv/bin/celery -A modules.campaign.celery_app:celery_app beat -l info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
